# templates/curator.cron.j2
0 1 * * * root {{curator_venv_path}}/bin/curator --config {{curator_venv_path}}/curator.yml {{curator_venv_path}}/curator_action_delete.yml
